# deploy/paper-codex.service
[Unit]
Description=Paper Codex research workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paper-codex
Group=paper-codex
WorkingDirectory=/opt/paper-codex
EnvironmentFile=/etc/paper-codex/paper-codex.env
ExecStart=/usr/local/bin/paper-codex
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
KillMode=mixed
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/paper-codex
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
UMask=0077

[Install]
WantedBy=multi-user.target
